# deploy/supervisor-app.conf
[program:postgresql-app]
command = postgresql -D /var/lib/postgresql/9.1/main -c config_file=/etc/postgresql/9.1/main/postgresql.conf

[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /opt/TimeSide/examples/sandbox/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx

