# examples/deploy/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /opt/TimeSide/examples/sandbox/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx

