# app/deploy/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /srv/app/sandbox/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
